# hugo.service
[Unit]
Description=Hugo Server

[Service]
User=hugo
Type=forking
ExecStart=/var/hugo/start.sh

[Install]
WantedBy=multi-user.target
