# hugo.service
[Unit]
Description=Hugo Server

[Service]
User=hugo
Type=forking
ExecStart=/home/hugo/hugo-config/start.sh

[Install]
WantedBy=multi-user.target
